# scripts/gramfix-hotkey.service.template
[Unit]
Description=GramFix Hotkey Daemon (Alt+G grammar correction)
Documentation=https://github.com/AniruthKarthik/gramfix
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=__HOTKEY_BIN__ --gramfix=__GRAMFIX_BIN__
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical-session.target
